# config/guide-api.service.example
[Unit]
Description=KKFlow Guide API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/kkflow-guide
EnvironmentFile=/opt/kkflow-guide/apps/guide-api/.env
ExecStart=/usr/bin/node apps/guide-api/src/server.js
Restart=always
RestartSec=3
User=www
Group=www

[Install]
WantedBy=multi-user.target
